# panopticon/panopticon.service
[Unit]
Description=Panopticon - RFID Door Access Web Server
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=panopticon
Group=panopticon
StateDirectory=panopticon
ExecStart=/usr/local/bin/panopticon
Environment=RUST_LOG=panopticon=debug,tower_http=info
EnvironmentFile=/etc/panopticon/env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
